# debian/wb-demo-kit-configs.postinst
#!/bin/bash

#DEBHELPER#

# restart services
invoke-rc.d wb-mqtt-serial restart
invoke-rc.d wb-rules restart
